# docker/crontab
0 15 * * *root /bin/bash /bilibili/job.sh >> /var/log/cron.log
